# servizi_demoni/sync_rtc.timer
[Unit]
Description=Sync RTC timer for Arduino

[Timer]
OnBootSec=1min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target


